# compass-harvester/crontab
PATH=/usr/local/bin
0 2 * * * /usr/bin/python /src/harvest_records.py /testfiles/ > /var/log/harvest-cron.log 2>&1
# Empty mandatory line
